# root/etc/s6-overlay/s6-rc.d/service-restic_rclone/run
#!/command/with-contenv bash

umask "${UMASK}"

exec s6-setuidgid restic /usr/local/bin/start.sh
